# sparse/lib/systemd/system/enable_modem.service
[Unit]
Description=Enable qmi modem

[Service]
Type=oneshot
ExecStart=/usr/bin/enable_modem.sh

[Install]
WantedBy=network.target
